# deploy/systemd/m-engine-api.service
# M-Engine — API FastAPI (uvicorn) como serviço systemd (runtime no HOST).
#
# Roda como o usuário cuja sessão Claude Code está logada (ex.: `ubuntu`), para que o
# provider "cc" reaproveite a auth do sistema (assinatura) sem crédito de API.
#
# Instalação (na VM):
#   # venv já criado: /home/ubuntu/m-venv  (pip install /opt/m-engine)
#   sudo cp deploy/systemd/m-engine-api.service /etc/systemd/system/
#   sudo install -m 0600 /etc/m-engine.env /etc/m-engine.env   # segredos 0600 (M_BASE, ELEVENLABS_API_KEY, M_FORCE_MODEL=cc, M_CLAUDE_CLI_BIN, M_API_HOST=<ip tailnet>, REDIS_URL)
#   sudo systemctl daemon-reload && sudo systemctl enable --now m-engine-api

[Unit]
Description=M-Engine API (FastAPI/uvicorn)
After=network-online.target redis-server.service tailscaled.service
Wants=network-online.target tailscaled.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu

# Segredos e config vêm exclusivamente do EnvironmentFile (modo 0600, fora do git).
EnvironmentFile=/etc/m-engine.env

WorkingDirectory=/opt/m-engine
# Bind no IP da tailnet (M_API_HOST) — API privada. Restart cobre a corrida com o tailscaled.
ExecStart=/home/ubuntu/m-venv/bin/uvicorn m_engine.api:app --host ${M_API_HOST} --port ${M_API_PORT}

Restart=always
RestartSec=3

# --- Hardening (PHI) ---
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
# ProtectHome=false: o venv, o binário `claude` e a auth (~/.claude) vivem em /home/ubuntu.
ProtectHome=false
# Único caminho gravável fora de /home: o volume de dados clínicos.
ReadWritePaths=/var/lib/m-data
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
UMask=0077

[Install]
WantedBy=multi-user.target
