# deploy/systemd/m-engine-api.service
# M-Engine — API FastAPI (uvicorn) como serviço systemd.
#
# Instalação (na VM de produção):
#   sudo useradd --system --home /opt/m-engine --shell /usr/sbin/nologin mengine
#   sudo mkdir -p /opt/m-engine && sudo chown mengine:mengine /opt/m-engine
#   # crie o venv e instale: sudo -u mengine python3.11 -m venv /opt/m-engine/venv
#   #                        sudo -u mengine /opt/m-engine/venv/bin/pip install /caminho/do/projeto
#   sudo cp deploy/systemd/m-engine-api.service /etc/systemd/system/
#   sudo install -m 0600 -o mengine -g mengine .env /etc/m-engine.env   # segredos: 0600
#   sudo systemctl daemon-reload && sudo systemctl enable --now m-engine-api

[Unit]
Description=M-Engine API (FastAPI/uvicorn)
After=network-online.target redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=mengine
Group=mengine

# Segredos e config vêm exclusivamente do EnvironmentFile (modo 0600, fora do git).
EnvironmentFile=/etc/m-engine.env

WorkingDirectory=/opt/m-engine
ExecStart=/opt/m-engine/venv/bin/uvicorn m_engine.api:app --host ${M_API_HOST} --port ${M_API_PORT}

Restart=always
RestartSec=3

# --- Hardening (PHI) ---
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
# Único caminho gravável: o volume de dados clínicos.
ReadWritePaths=/var/lib/m-engine
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
UMask=0077

[Install]
WantedBy=multi-user.target
