# deploy/systemd/m-engine-worker.service
# M-Engine — worker Celery como serviço systemd (runtime no HOST).
#
# Pré-requisitos: mesmo usuário/venv/EnvironmentFile da API (ver m-engine-api.service)
# e um Redis local (REDIS_URL=redis://localhost:6379/0). Roda como `ubuntu` para que o
# provider "cc" use a auth da assinatura (Claude Code) sem crédito de API.
#
# Instalação:
#   sudo cp deploy/systemd/m-engine-worker.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now m-engine-worker

[Unit]
Description=M-Engine Worker (Celery)
After=network-online.target redis-server.service
Wants=network-online.target
Requires=redis-server.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu

# Mesmas variáveis da API (ELEVENLABS_API_KEY, REDIS_URL, M_BASE, M_FORCE_MODEL, M_CLAUDE_CLI_BIN).
EnvironmentFile=/etc/m-engine.env

WorkingDirectory=/opt/m-engine
ExecStart=/home/ubuntu/m-venv/bin/celery -A m_engine.tasks worker --loglevel=info --concurrency=2

Restart=always
RestartSec=3
# Encerramento gracioso: deixa a task em andamento terminar antes do SIGKILL.
TimeoutStopSec=120
KillSignal=SIGTERM

# --- Hardening (PHI) ---
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
# ProtectHome=false: venv + binário `claude` + auth (~/.claude) vivem em /home/ubuntu.
ProtectHome=false
ReadWritePaths=/var/lib/m-data
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
UMask=0077

[Install]
WantedBy=multi-user.target
